# files/etc/systemd/system/auto-hibernate.service
[Unit]
Description=Auto hibernate if battery is low

[Service]
Type=simple
ExecStart=/usr/local/bin/auto-hibernate
