# files/etc/systemd/system/detect-egpu.service
[Unit]
Description=Toggle eGPU support

[Service]
ExecStart=/usr/bin/detect-egpu
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
