# files/etc/systemd/system/pacman-index.service
[Unit]
Description=Update the pacman package index

[Service]
Type=oneshot
ExecStart=/usr/bin/chronic pacman -Suwy --noconfirm
