# files/etc/systemd/system/powertop.service
[Unit]
Description=Power usage optimisations

[Service]
ExecStart=/usr/bin/powertop-safe
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
