# sysconfig/etc/systemd/system/pacman-refresh.service
[Unit]
Description=Update pacman index and retrieve packages

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman -Suwy --noconfirm
