# common/.config/systemd/user/assia-vpn.service
[Unit]
Description=Assia Cisco VPN

[Service]
ExecStart=/bin/sh %h/.bin/assia-vpn

[Install]
WantedBy=desktop.target
