# common/.config/systemd/user/blueman-applet.service
[Unit]
Description=A tray applet for managing bluetooth

[Service]
ExecStart=/usr/bin/blueman-applet

[Install]
WantedBy=desktop.target
