# common/.config/systemd/user/blueman-applet.service
[Unit]
Description=A tray applet for managing bluetooth
After=xorg.target

[Service]
ExecStart=/usr/bin/blueman-applet

[Install]
WantedBy=wm.target
