# common/.config/systemd/user/compiz.service
[Unit]
Description=Composite window manager
Before=compizDesktop.target

[Service]
ExecStart=/usr/bin/compiz ccp
Restart=always
RestartSec=7

[Install]
WantedBy=desktop.target
