# common/.config/systemd/user/davmail.service
[Unit]
Description=A POP/IMAP/SMTP/Caldav/Carddav/LDAP exchange gateway.

[Service]
ExecStart=/usr/bin/davmail

[Install]
WantedBy=desktop.target
