# common/.config/systemd/user/desktop.target
[Unit]
Description=Compiz-based desktop
Wants=wm.target

[Install]
Alias=default.target
