# common/.config/systemd/user/fehbg.service
[Unit]
Description=An wallpaper setter

[Service]
ExecStart=/bin/sh /home/hugo/.fehbg
Type=oneshot

[Install]
WantedBy=desktop.target
