# common/.config/systemd/user/fehbg.service
[Unit]
Description=An wallpaper setter
After=xorg.target

[Service]
ExecStart=/bin/sh $HOME/.fehbg
Type=oneshot

[Install]
WantedBy=wm.target
