# common/.config/systemd/user/lxpanel.service
[Unit]
Description=A lightweight GTK2-based panel

[Service]
ExecStart=/usr/bin/lxpanel

[Install]
WantedBy=desktop.target
