# common/.config/systemd/user/numlock.service
[Unit]
Description=Activate the numeric keypad

[Service]
ExecStart=/usr/bin/numlockx on
Type=oneshot

[Install]
WantedBy=desktop.target
