# common/.config/systemd/user/parcellite.service
[Unit]
Description=Lightweight GTK+ Clipboard Manager

[Service]
ExecStart=/usr/bin/parcellite --daemon --no-icon

[Install]
WantedBy=desktop.target
