# common/.config/systemd/user/sxhkd.service
[Unit]
Description=A grabbing keys program for X

[Service]
ExecStart=/usr/bin/sxhkd
KillMode=process

[Install]
WantedBy=desktop.target
