# common/.config/systemd/user/thunar.service
[Unit]
Description=File Manager

[Service]
ExecStart=/usr/bin/thunar --daemon

[Install]
WantedBy=desktop.target
