# common/.config/systemd/user/thunar.service
[Unit]
Description=File Manager
After=xorg.target

[Service]
ExecStart=/usr/bin/thunar --daemon

[Install]
WantedBy=wm.target
