# common/.config/systemd/user/urxvtd.service
[Unit]
Description=urxvt terminal daemon

[Service]
ExecStart=/usr/bin/urxvtd

[Install]
WantedBy=desktop.target
