# common/.config/systemd/user/volctld.service
[Unit]
Description=A volume control daemon
After=xorg.target

[Service]
ExecStart=/usr/bin/volctld

[Install]
WantedBy=wm.target
