# common/.config/systemd/user/wicd-gtk.service
[Unit]
Description=GTK client for wicd
After=xorg.target

[Service]
ExecStart=/usr/bin/wicd-gtk -t

[Install]
WantedBy=wm.target
