# common/.config/systemd/user/wicd-gtk.service
[Unit]
Description=GTK client for wicd

[Service]
ExecStart=/usr/bin/wicd-gtk -t

[Install]
WantedBy=desktop.target
