# common/.config/systemd/user/xbindkeys.service
[Unit]
Description=A grabbing keys program for X
After=xorg.target

[Service]
ExecStart=/usr/bin/xbindkeys -n

[Install]
WantedBy=wm.target
