# home/.config/cron/root.crontab
0	*/6	*	*	*	/usr/bin/pacman -Suwy --noconfirm > /dev/null
