# etc/systemd/system/getty@tty1.service.d/slow-restart.conf
# Services are restarted in 100ms by default.
# That's a bit too fast for this service, since it restart sway and the entire
# graphical session.
[Service]
RestartSec=5s
