# common/.config/cron/crontab
# Restart offlineimap on an hourly basis since it tends to stall every
# once in a while.
0	*	*	*	*	/bin/sh /home/hugo/.config/cron/offlineimap.cron
