# common/.config/cron/crontab
# Restart offlineimap on an hourly basis since it tends to stall every
# once in a while.
0	*	*	*	*	/usr/bin/usrc restart offlineimap > /dev/null
