# common/.config/cron/hugo.cron
# Restart offlineimap on an hourly basis since it tends to stall every
# once in a while.
0	*	*	*	*	/usr/bin/systemctl --user restart offlineimap
