# home/.config/systemd/user/dump-packages.service
[Unit]
Description=Dump currently installed packages

[Service]
ExecStart=/bin/sh %h/.local/lib/dump-packages
Type=simple
