# home/.config/systemd/user/i3.service
[Unit]
Description=A tiling window manager
After=xorg.target
Requires=xorg.target

[Service]
ExecStart=/usr/bin/i3

[Install]
WantedBy=xorg.target
