# home/.config/systemd/user/khal.service
[Unit]
Description=Update khal's cache

[Service]
ExecStart=/usr/bin/khal
Type=oneshot
StandardOutput=null
