# home/.config/systemd/user/khal.service
[Unit]
Description=Update khal's cache
Documentation=man:khal(1)

[Service]
ExecStart=/usr/bin/khal
Type=oneshot
StandardOutput=null
