# home/.config/systemd/user/notmuch.service
[Unit]
Description=Compact Maildirs
Documentation=man:notmuch(1)

[Service]
ExecStart=/usr/bin/notmuch compact
Type=simple
