# home/.config/systemd/user/notmuch.service
[Unit]
Description=Compact Maildirs

[Service]
ExecStart=/usr/bin/notmuch compact
Type=simple

[Install]
WantedBy=default.target
