# home/.config/systemd/user/numlock.service
[Unit]
After=xorg.target
Description=Enable numlock

[Service]
Type=oneshot
ExecStart=/usr/bin/numlockx

[Install]
WantedBy=xorg.target
