# home/.config/systemd/user/offlineimap.service
[Unit]
Description=Synchronize Maildirs
Documentation=man:offlineimap(1)

[Service]
ExecStart=/usr/bin/offlineimap
Type=simple
Restart=on-abnormal
