# home/.config/systemd/user/offlineimap.service
[Unit]
Description=Synchronize Maildirs

[Service]
ExecStart=/usr/bin/offlineimap
Type=simple
Restart=on-abnormal

[Install]
WantedBy=default.target
