# home/.config/systemd/user/offlineimap.timer
[Unit]
Description=Synchronize Maildirs

[Timer]
OnUnitInactiveSec=5m

[Install]
WantedBy=timers.target
