# home/.config/systemd/user/status-email@.service
[Unit]
Description=Send summary email for %i

[Service]
ExecStart=/usr/bin/sh %h/.local/bin/systemd-email %i
Type=oneshot
