# home/.config/systemd/user/vdirsyncer.service
[Unit]
Description=Synchronize vdirs

[Service]
ExecStart=/usr/bin/vdirsyncer -v DEBUG sync
ExecStartPost=khal
Type=simple
Restart=on-abnormal
