# home/.config/systemd/user/vdirsyncer.service
[Unit]
Description=Synchronize vdirs
Requires=khal.service
Before=khal.service

[Service]
ExecStart=/usr/bin/vdirsyncer sync
Type=oneshot
