# home/.config/systemd/user/vdirsyncer.timer
[Unit]
Description=Synchronize vdirs

[Timer]
OnBootSec=5m
OnUnitActiveSec=10m
AccuracySec=5m

[Install]
WantedBy=timers.target
